# crontab
0 0 * * * root find /jobs -type d -mtime +30 -exec /SEACR-webapp/clean.sh {} \; >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job
